# project-base/app/supervisor/dev/worker/product_reindex.conf
[program:product_reindex]
command=php /usr/local/var/www/shopsys/bin/console rabbitmq:consumer product_reindex --messages=250 --env=dev
process_name=%(program_name)s%(process_num)02d
numprocs=1
startsecs=2
autorestart=true
stopsignal=INT
stopasgroup=true
stopwaitsecs=60
stdout_logfile=/usr/local/var/www/shopsys/project-base/app/supervisor/dev/logs/stdout.log
stderr_logfile=/usr/local/var/www/shopsys/project-base/app/supervisor/dev/logs/stderr.log

